# templates/default/ora_init_script.erb
#!/bin/sh
#
# Init file for Oracle database server
#
# chkconfig: 35 99 02
# description: Oracle database server
 
case "$1" in
'start')
        su - oracle -c "<%= node[:oracle][:rdbms][:ora_home] %>/bin/dbstart <%= node[:oracle][:rdbms][:ora_home] %>"
        su - oracle -c "<%= node[:oracle][:rdbms][:ora_home] %>/bin/emctl start dbconsole"
        RETVAL=$?
        [ "$RETVAL" = 0 ] && touch /var/lock/subsys/oracle
        ;;
 
'stop')
        su - oracle -c "<%= node[:oracle][:rdbms][:ora_home] %>/bin/emctl stop dbconsole"
        su - oracle -c "<%= node[:oracle][:rdbms][:ora_home] %>/bin/dbshut <%= node[:oracle][:rdbms][:ora_home] %>"
        RETVAL=$?
        [ "$RETVAL" = 0 ] && rm -f /var/lock/subsys/oracle
        ;;
 
*)
        echo "Usage: $0 { start | stop }"
        exit 1
        ;;
 
esac
exit 0
